# files/php/cronjob-worker.conf
[program:cronjob-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/cron
user=root
